# core/defaults/upstart/moonstalk.conf
start on (local-filesystems and net-device-up) or moonstalk
stop on runlevel [0] or killmoonstalk

task 
chdir /usr/local/moonstalk
exec ./runner start all >> /tmp/moonstalk-upstart.log 2>&1
