# core/defaults/upstart/moonstalk-stop.conf

start on runlevel [0] or stopmoonstalk

task 
chdir /usr/local/moonstalk
exec ./runner stop all >> /tmp/moonstalk-upstart.log 2>&1
